# compute/alpine/services/api
#!/sbin/openrc-run
# Opentrons API server start script

depend() {
  need localmount
}

start() {
  local retval

  export ENABLE_VIRTUAL_SMOOTHIE=true

  ebegin "Starting Opentrons API server"
  start-stop-daemon --start \
      --background \
      --pidfile /var/run/api.pid \
      --make-pidfile \
      --startas /bin/bash \
      -- -c "python -m opentrons.server.main 0.0.0.0:31950 > /var/log/api.log 2>&1"
  retval=$?
  eend ${retval}
  return ${retval}
}

stop() {
  local retval

  ebegin "Stopping Opentrons API server"
  start-stop-daemon --stop --pidfile /var/run/api.pid
  retval=$?
  eend ${retval}
  return ${retval}
}