# compute/alpine/services/mdns-advertiser
#!/sbin/openrc-run
# mdns service advertiser start script

depend() {
  need localmount
}

start() {
  local retval
  export DBUS_SYSTEM_BUS_ADDRESS=unix:path=/host/run/dbus/system_bus_socket
  ebegin "Starting mdns service advertiser"
  start-stop-daemon --start \
      --background \
      --pidfile /var/run/mdns-advertiser.pid \
      --make-pidfile \
      --exec announce_mdns.py
  retval=$?
  eend ${retval}
  return ${retval}
}

stop() {
  local retval

  ebegin "Stopping mdns service advertiser"
  start-stop-daemon --stop --pidfile /var/run/mdns-advertiser.pid
  retval=$?
  eend ${retval}
  return ${retval}
}